# board/yarmp/overlay/etc/inittab
# Startup the system
null::sysinit:/sbin/modprobe -q brcmfmac
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mkdir -p /dev/pts /dev/shm
null::sysinit:/bin/mount -a
null::sysinit:/bin/ln -sf /proc/self/fd /dev/fd
null::sysinit:/bin/ln -sf /proc/self/fd/0 /dev/stdin
null::sysinit:/bin/ln -sf /proc/self/fd/1 /dev/stdout
null::sysinit:/bin/ln -sf /proc/self/fd/2 /dev/stderr
null::sysinit:/bin/hostname -F /etc/hostname
null::sysinit:/bin/chown :555 /dev/input/event* /dev/ttyAMA0 /dev/snd/*
null::sysinit:/bin/chmod g+rw /dev/input/event* /dev/ttyAMA0 /dev/snd/*
null::sysinit:/sbin/ifup lo

#
# start yarmp
#


#
#
#

null::sysinit:/sbin/ifup wlan0
null::sysinit:/sbin/dropbear -jk -b /etc/issue
#ttyAMA0::respawn:/sbin/getty -L ttyAMA0 115200 vt100 # GENERIC_SERIAL

# Stuff to do before rebooting
null::shutdown:/sbin/ifdown -a
null::shutdown:/bin/umount -a -r
