# board/yarmp/overlay/etc/profile.d/alias.sh
alias ll='ls -la --color=always'
alias vim='vi'
